# deploy/maintainer.service
[Unit]
Description=Maintainer host-maintenance daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=maintainer
Group=maintainer
WorkingDirectory=/opt/maintainer
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 main.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
