# deploy/maintainer.service
[Unit]
Description=Maintainer host-maintenance daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=simon
Group=simon
WorkingDirectory=/home/simon/maintainer
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/simon/maintainer/.venv/bin/python3 /home/simon/maintainer/main.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
